# luci-app-singbox-ui/root/etc/uci-defaults/99-luci-app-singbox-ui
#!/bin/sh
# Run by OpenWrt on first boot after package install. Removes itself on success.

# Force LuCI's RPC ACL cache to reload so the menu entry appears immediately.
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null

# Reload rpcd so the new ACL JSON is picked up.
/etc/init.d/rpcd reload >/dev/null 2>&1

# Install cron entry: refresh subscriptions+rulesets every 30 minutes.
# Per-section update_interval / sub_interval is honoured by subscription.uc itself.
CRON_LINE='*/30 * * * * /usr/bin/ucode /usr/share/singbox-ui/subscription.uc refresh all >/dev/null 2>&1'
CRON_FILE=/etc/crontabs/root
mkdir -p "$(dirname "$CRON_FILE")"
touch "$CRON_FILE"
# Drop legacy refresh.sh line from prior package versions, then ensure the
# current ucode line is present exactly once.
if grep -qF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
	|| ! grep -qF '/usr/share/singbox-ui/subscription.uc' "$CRON_FILE"; then
	tmp=$(mktemp)
	grep -vF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
		| grep -vF '/usr/share/singbox-ui/subscription.uc' > "$tmp"
	echo "$CRON_LINE" >> "$tmp"
	mv "$tmp" "$CRON_FILE"
	/etc/init.d/cron restart >/dev/null 2>&1 || /etc/init.d/cron start >/dev/null 2>&1
fi

exit 0
