# luci-app-singbox-ui/root/etc/uci-defaults/99-luci-app-singbox-ui
#!/bin/sh
# Run by OpenWrt on first boot after package install. Removes itself on success.

# Force LuCI's RPC ACL cache to reload so the menu entry appears immediately.
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null

# Reload rpcd so the new ACL JSON is picked up.
/etc/init.d/rpcd reload >/dev/null 2>&1

# Install cron entry: refresh subscriptions+rulesets every 30 minutes.
# Per-section update_interval / sub_interval is honoured by subscription.uc itself.
CRON_LINE='*/30 * * * * /usr/bin/ucode -L /usr/share/singbox-ui/lib /usr/share/singbox-ui/subscription.uc refresh all >/dev/null 2>&1'
CRON_FILE=/etc/crontabs/root
mkdir -p "$(dirname "$CRON_FILE")"
touch "$CRON_FILE"
# Drop legacy refresh.sh line from prior package versions, then ensure the
# current ucode line is present exactly once.
if grep -qF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
	|| ! grep -qF "$CRON_LINE" "$CRON_FILE"; then
	tmp=$(mktemp)
	grep -vF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
		| grep -vF '/usr/share/singbox-ui/subscription.uc' > "$tmp"
	echo "$CRON_LINE" >> "$tmp"
	mv "$tmp" "$CRON_FILE"
	/etc/init.d/cron restart >/dev/null 2>&1 || /etc/init.d/cron start >/dev/null 2>&1
fi

# ---------------------------------------------------------------------------
# Migrate legacy FakeIP list form (multiple `list inet*_range`) to scalar
# `option inet*_range`. sing-box 1.12+ rejects array form; we now expose only
# the first element in the UI. Other elements are dropped with a warning.
# ---------------------------------------------------------------------------
migrate_fakeip_list_to_option() {
	# $1 = inet4_range | inet6_range
	# `uci get` for a list option returns values space-joined.
	_mig_raw=$(uci -q get singbox-ui.fakeip."$1" 2>/dev/null || true)
	[ -n "$_mig_raw" ] || return 0

	# If it was already a scalar option, `uci show` won't include a numeric
	# index. List form shows like:  singbox-ui.fakeip.inet4_range='1.2.3.0/24 4.5.6.0/16'
	# in `uci get`, and `uci show` lists multiple `=` lines per element only
	# in newer uci. The simplest heuristic: count words in $_mig_raw.
	_mig_first="${_mig_raw%% *}"
	_mig_rest="${_mig_raw#"$_mig_first"}"
	if [ -n "$_mig_rest" ]; then
		logger -t singbox-ui "FakeIP $1 migration: keeping '$_mig_first', dropping: $_mig_rest"
	fi

	# Re-set as a scalar option (uci delete first removes the list typing).
	uci -q delete singbox-ui.fakeip."$1" >/dev/null 2>&1 || true
	uci -q set    singbox-ui.fakeip."$1"="$_mig_first"
	uci -q commit singbox-ui
}

migrate_fakeip_list_to_option inet4_range
migrate_fakeip_list_to_option inet6_range

# ---------------------------------------------------------------------------
# Migrate the legacy single `tproxy` section into an `inbound` section
# (protocol=tproxy). The Inbounds tab replaces the old derived model.
# Idempotent: only runs while a `tproxy` section still exists, and removes it.
# ---------------------------------------------------------------------------
migrate_tproxy_to_inbound() {
	# Section type check: `uci get singbox-ui.tproxy` returns the type name.
	_t=$(uci -q get singbox-ui.tproxy 2>/dev/null || true)
	[ "$_t" = "tproxy" ] || return 0

	_en=$(uci -q get singbox-ui.tproxy.enabled 2>/dev/null || echo 0)
	_port=$(uci -q get singbox-ui.tproxy.port 2>/dev/null || echo 7893)
	_hj=$(uci -q get singbox-ui.tproxy.hijack_dns 2>/dev/null || echo 0)
	_ifaces=$(uci -q get singbox-ui.tproxy.interface 2>/dev/null || echo br-lan)

	# Create (or overwrite) the named inbound section.
	uci -q set singbox-ui.tproxy_in='inbound'
	uci -q set singbox-ui.tproxy_in.enabled="$_en"
	uci -q set singbox-ui.tproxy_in.mode='constructor'
	uci -q set singbox-ui.tproxy_in.protocol='tproxy'
	uci -q set singbox-ui.tproxy_in.listen='::'
	uci -q set singbox-ui.tproxy_in.listen_port="$_port"
	uci -q set singbox-ui.tproxy_in.hijack_dns="$_hj"
	uci -q set singbox-ui.tproxy_in.nft_rules='1'
	uci -q delete singbox-ui.tproxy_in.interface 2>/dev/null || true
	for w in $_ifaces; do
		uci -q add_list singbox-ui.tproxy_in.interface="$w"
	done

	uci -q delete singbox-ui.tproxy
	uci -q commit singbox-ui
}

migrate_tproxy_to_inbound

# ---------------------------------------------------------------------------
# Drop the removed per-outbound `expose_*` options from all outbound sections.
# ---------------------------------------------------------------------------
migrate_drop_expose() {
	_changed=0
	for _sec in $(uci -q show singbox-ui 2>/dev/null \
		| sed -n 's/^singbox-ui\.\([^.]*\)=outbound$/\1/p'); do
		for _o in expose_proxy expose_type expose_port expose_listen; do
			if uci -q get "singbox-ui.$_sec.$_o" >/dev/null 2>&1; then
				uci -q delete "singbox-ui.$_sec.$_o"
				_changed=1
			fi
		done
	done
	[ "$_changed" = "1" ] && uci -q commit singbox-ui
	return 0
}

migrate_drop_expose

# ---------------------------------------------------------------------------
# Create default sections for Phase 2 features if they don't exist yet.
# These match `etc/config/singbox-ui` defaults and let users see the sections
# in LuCI without manually editing UCI.
# ---------------------------------------------------------------------------
ensure_section() {
	# $1 = section name, $2 = section type
	# Returns 0 if newly created, 1 if it already existed (caller can branch).
	if ! uci -q get "singbox-ui.$1" >/dev/null 2>&1; then
		uci -q set "singbox-ui.$1=$2"
		return 0
	fi
	return 1
}

if ensure_section log log; then
	uci -q set singbox-ui.log.enabled='1'
	uci -q set singbox-ui.log.level='info'
	uci -q commit singbox-ui
fi

if ensure_section cache cache; then
	uci -q set singbox-ui.cache.enabled='0'
	uci -q set singbox-ui.cache.path='/tmp/singbox-ui-cache.db'
	uci -q commit singbox-ui
fi

# ---------------------------------------------------------------------------
# DNS model migration to typed dns_server/dns_rule/dns (piece C).
# fakeip list→option migration above runs first, so ranges are scalar here.
# ---------------------------------------------------------------------------
migrate_fakeip_to_dns_server() {
	[ "$(uci -q get singbox-ui.fakeip 2>/dev/null)" = "fakeip" ] || return 0
	_en=$(uci -q get singbox-ui.fakeip.enabled 2>/dev/null || echo 0)
	_v4=$(uci -q get singbox-ui.fakeip.inet4_range 2>/dev/null || echo '198.18.0.0/15')
	_v6=$(uci -q get singbox-ui.fakeip.inet6_range 2>/dev/null || echo 'fc00::/18')
	uci -q delete singbox-ui.fakeip
	uci -q set singbox-ui.fakeip='dns_server'
	uci -q set singbox-ui.fakeip.enabled="$_en"
	uci -q set singbox-ui.fakeip.type='fakeip'
	uci -q set singbox-ui.fakeip.inet4_range="$_v4"
	[ -n "$_v6" ] && uci -q set singbox-ui.fakeip.inet6_range="$_v6"
	uci -q commit singbox-ui
}
migrate_fakeip_to_dns_server

migrate_dns_outbound_to_server() {
	[ "$(uci -q get singbox-ui.dns_outbound 2>/dev/null)" = "dns_outbound" ] || return 0
	_addr=$(uci -q get singbox-ui.dns_outbound.address 2>/dev/null || echo '')
	_det=$(uci -q get singbox-ui.dns_outbound.detour 2>/dev/null || echo 'direct')
	uci -q set singbox-ui.out_dns='dns_server'
	uci -q set singbox-ui.out_dns.enabled='1'
	uci -q set singbox-ui.out_dns.detour="$_det"
	case "$_addr" in
		https://*) uci -q set singbox-ui.out_dns.type='https'
		           _hp=${_addr#https://}; _host=${_hp%%/*}; _path=/${_hp#*/}
		           uci -q set singbox-ui.out_dns.server="${_host%%:*}"
		           uci -q set singbox-ui.out_dns.server_port='443'
		           [ "$_path" != "/$_hp" ] && uci -q set singbox-ui.out_dns.path="$_path" ;;
		tls://*)   uci -q set singbox-ui.out_dns.type='tls'
		           uci -q set singbox-ui.out_dns.server="${_addr#tls://}" ;;
		*)         uci -q set singbox-ui.out_dns.type='udp'
		           uci -q set singbox-ui.out_dns.server="$_addr" ;;
	esac
	uci -q set singbox-ui.dns='dns' 2>/dev/null || true
	uci -q set singbox-ui.dns.final='out_dns'
	uci -q delete singbox-ui.dns_outbound
	uci -q commit singbox-ui
}
migrate_dns_outbound_to_server

migrate_ruleset_dns_fakeip() {
	_changed=0
	for _sec in $(uci -q show singbox-ui 2>/dev/null | sed -n 's/^singbox-ui\.\([^.]*\)=ruleset$/\1/p'); do
		if [ "$(uci -q get "singbox-ui.$_sec.dns_fakeip" 2>/dev/null)" = "1" ]; then
			_tag=$(uci -q get "singbox-ui.$_sec.dns_fakeip_tag" 2>/dev/null || echo 'fakeip')
			_rule="dnsrule_$_sec"
			uci -q set "singbox-ui.$_rule=dns_rule"
			uci -q set "singbox-ui.$_rule.enabled=1"
			uci -q delete "singbox-ui.$_rule.ruleset" 2>/dev/null || true
			uci -q add_list "singbox-ui.$_rule.ruleset=$_sec"
			uci -q set "singbox-ui.$_rule.server=$_tag"
			_changed=1
		fi
		for _o in dns_fakeip dns_fakeip_tag; do
			if uci -q get "singbox-ui.$_sec.$_o" >/dev/null 2>&1; then
				uci -q delete "singbox-ui.$_sec.$_o"; _changed=1
			fi
		done
	done
	[ "$_changed" = "1" ] && uci -q commit singbox-ui
	return 0
}
migrate_ruleset_dns_fakeip

# ---------------------------------------------------------------------------
# Ensure a clash_api section exists and has a stable random secret.
# ---------------------------------------------------------------------------
if ensure_section clash_api clash_api; then
	uci -q set singbox-ui.clash_api.enabled='0'
	uci -q set singbox-ui.clash_api.listen='127.0.0.1'
	uci -q set singbox-ui.clash_api.port='9090'
	uci -q commit singbox-ui
fi

ensure_clash_api_secret() {
	[ "$(uci -q get singbox-ui.clash_api 2>/dev/null)" = "clash_api" ] || return 0
	_sec=$(uci -q get singbox-ui.clash_api.secret 2>/dev/null || true)
	[ -n "$_sec" ] && return 0
	_new=$(head -c16 /dev/urandom 2>/dev/null | hexdump -v -e '/1 "%02x"' 2>/dev/null)
	[ -n "$_new" ] || _new=$(date +%s)$$
	uci -q set singbox-ui.clash_api.secret="$_new"
	uci -q commit singbox-ui
}
ensure_clash_api_secret

# ---------------------------------------------------------------------------
# Cache file migration: convert legacy `cache` (often disabled + explicit
# /tmp path) into the new storage-mode model. Idempotent.
# ---------------------------------------------------------------------------
migrate_cache_storage() {
	[ "$(uci -q get singbox-ui.cache 2>/dev/null)" = "cache" ] || \
		uci -q set singbox-ui.cache='cache'

	# If storage is already set, the user (or this script) has migrated.
	[ -n "$(uci -q get singbox-ui.cache.storage 2>/dev/null)" ] && return 0

	_was_legacy=0
	_path=$(uci -q get singbox-ui.cache.path 2>/dev/null || true)
	case "$_path" in
		""|"/tmp/singbox-ui-cache.db")
			uci -q set singbox-ui.cache.storage='ram'
			uci -q delete singbox-ui.cache.path 2>/dev/null || true
			_was_legacy=1 ;;
		"/etc/sing-box/cache.db")
			uci -q set singbox-ui.cache.storage='flash'
			uci -q delete singbox-ui.cache.path 2>/dev/null || true ;;
		*)
			uci -q set singbox-ui.cache.storage='custom' ;;
	esac

	# Only flip the shipped-default disabled state. A user who explicitly
	# disabled cache after customising path keeps their setting.
	if [ "$_was_legacy" = "1" ]; then
		uci -q set singbox-ui.cache.enabled='1'
	fi

	# store_fakeip default: only set if absent.
	_sf=$(uci -q get singbox-ui.cache.store_fakeip 2>/dev/null || echo '')
	[ "$_sf" = "" ] && uci -q set singbox-ui.cache.store_fakeip='1'

	uci -q commit singbox-ui
}
migrate_cache_storage

# ---------------------------------------------------------------------------
# Create the default DNS inbound for upgrades. Idempotent: skip if a
# section with that name already exists (user may have customised it).
# ---------------------------------------------------------------------------
ensure_dns_inbound() {
	[ -n "$(uci -q get singbox-ui.dns_in 2>/dev/null)" ] && return 0
	uci -q set singbox-ui.dns_in='inbound'
	uci -q set singbox-ui.dns_in.enabled='1'
	uci -q set singbox-ui.dns_in.protocol='direct'
	uci -q set singbox-ui.dns_in.listen='127.0.0.53'
	uci -q set singbox-ui.dns_in.listen_port='53'
	uci -q set singbox-ui.dns_in.network='udp'
	uci -q set singbox-ui.dns_in.dns_listener='1'
	uci -q commit singbox-ui
}
ensure_dns_inbound

exit 0
