# luci-app-singbox-ui/root/etc/uci-defaults/99-luci-app-singbox-ui
#!/bin/sh
# Run by OpenWrt on first boot after package install. Removes itself on success.

# Force LuCI's RPC ACL cache to reload so the menu entry appears immediately.
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null

# Reload rpcd so the new ACL JSON is picked up.
/etc/init.d/rpcd reload >/dev/null 2>&1

# Install cron entry: refresh subscriptions+rulesets every 30 minutes.
# Per-section update_interval / sub_interval is honoured by subscription.uc itself.
CRON_LINE='*/30 * * * * /usr/bin/ucode -L /usr/share/singbox-ui/lib /usr/share/singbox-ui/subscription.uc refresh all >/dev/null 2>&1'
CRON_FILE=/etc/crontabs/root
mkdir -p "$(dirname "$CRON_FILE")"
touch "$CRON_FILE"
# Drop legacy refresh.sh line from prior package versions, then ensure the
# current ucode line is present exactly once.
if grep -qF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
	|| ! grep -qF "$CRON_LINE" "$CRON_FILE"; then
	tmp=$(mktemp)
	grep -vF '/usr/share/singbox-ui/refresh.sh' "$CRON_FILE" \
		| grep -vF '/usr/share/singbox-ui/subscription.uc' > "$tmp"
	echo "$CRON_LINE" >> "$tmp"
	mv "$tmp" "$CRON_FILE"
	/etc/init.d/cron restart >/dev/null 2>&1 || /etc/init.d/cron start >/dev/null 2>&1
fi

# ---------------------------------------------------------------------------
# Migrate legacy FakeIP list form (multiple `list inet*_range`) to scalar
# `option inet*_range`. sing-box 1.12+ rejects array form; we now expose only
# the first element in the UI. Other elements are dropped with a warning.
# ---------------------------------------------------------------------------
migrate_fakeip_list_to_option() {
	# $1 = inet4_range | inet6_range
	# `uci get` for a list option returns values space-joined.
	_mig_raw=$(uci -q get singbox-ui.fakeip."$1" 2>/dev/null || true)
	[ -n "$_mig_raw" ] || return 0

	# If it was already a scalar option, `uci show` won't include a numeric
	# index. List form shows like:  singbox-ui.fakeip.inet4_range='1.2.3.0/24 4.5.6.0/16'
	# in `uci get`, and `uci show` lists multiple `=` lines per element only
	# in newer uci. The simplest heuristic: count words in $_mig_raw.
	_mig_first="${_mig_raw%% *}"
	_mig_rest="${_mig_raw#"$_mig_first"}"
	if [ -n "$_mig_rest" ]; then
		logger -t singbox-ui "FakeIP $1 migration: keeping '$_mig_first', dropping: $_mig_rest"
	fi

	# Re-set as a scalar option (uci delete first removes the list typing).
	uci -q delete singbox-ui.fakeip."$1" >/dev/null 2>&1 || true
	uci -q set    singbox-ui.fakeip."$1"="$_mig_first"
	uci -q commit singbox-ui
}

migrate_fakeip_list_to_option inet4_range
migrate_fakeip_list_to_option inet6_range

# ---------------------------------------------------------------------------
# Migrate tproxy.interface from scalar `option` to `list`. nftables.uc now
# supports multiple input interfaces; list form is required.
# Idempotent: deleting + re-adding produces an identical state when already a list.
# ---------------------------------------------------------------------------
migrate_tproxy_interface_to_list() {
	_val=$(uci -q get singbox-ui.tproxy.interface 2>/dev/null || true)
	if [ -n "$_val" ]; then
		uci -q delete singbox-ui.tproxy.interface
		for w in $_val; do
			uci -q add_list singbox-ui.tproxy.interface="$w"
		done
		uci -q commit singbox-ui
	fi
}

migrate_tproxy_interface_to_list

exit 0
