# luci-singbox-ui/root/etc/init.d/singbox-ui
#!/bin/sh /etc/rc.common
# shellcheck disable=SC2034
USE_PROCD=1
START=99
STOP=01

GENERATE_UC=/usr/share/singbox-ui/generate.uc
SUBSCRIPTION_UC=/usr/share/singbox-ui/subscription.uc
NFTABLES_UC=/usr/share/singbox-ui/nftables.uc
UCODE_LIB=/usr/share/singbox-ui/lib
# Overridable so tests can stub the binary (default is the packaged path used by
# the procd command below). Used for both `check` (S6.1) and `run`.
SINGBOX_BIN="${SINGBOX_BIN:-/usr/bin/sing-box}"

start_service() {
	mkdir -p /tmp/singbox-ui

	# Always try to refresh on (re)start, with short curl timeouts so boot
	# stays snappy. subscription.uc keeps prior sub_*.txt / rs_*.json files
	# untouched on failure, so generate.uc falls back to the cache.
	SINGBOX_BOOT_FETCH=1 ucode -L "$UCODE_LIB" "$SUBSCRIPTION_UC" fetch-subs &
	SINGBOX_BOOT_FETCH=1 ucode -L "$UCODE_LIB" "$SUBSCRIPTION_UC" fetch-rulesets &
	wait

	ucode -L "$UCODE_LIB" "$GENERATE_UC"

	# If neither a fresh download nor a cached config produced anything,
	# refuse to start; cron will re-run subscription.uc in ≤30 min and call
	# `reload`, which triggers start_service again.
	if [ ! -s /tmp/singbox-ui.json ]; then
		logger -t singbox-ui "refusing to start: no cached config and fetch failed"
		return 1
	fi

	# S6.1: validate the config semantically before handing it to procd. A
	# well-formed-but-rejected config (unknown field, dangling reference, bad
	# value) would otherwise make procd respawn sing-box in a tight loop. Refuse
	# to start instead and let cron/reload retry in ≤30 min. S1.1: record the
	# real apply outcome in apply_state.json (distinct from generate's state).
	mkdir -p /var/lib/singbox-ui
	if ! "$SINGBOX_BIN" check -c /tmp/singbox-ui.json 2>/tmp/singbox-ui.check.err; then
		logger -t singbox-ui "refusing to start: sing-box rejected config"
		logger -t singbox-ui "sing-box check: $(head -c 300 /tmp/singbox-ui.check.err 2>/dev/null)"
		echo "{\"last_apply_ts\":$(date +%s),\"last_apply_result\":\"failed\"}" \
			> /var/lib/singbox-ui/apply_state.json
		return 1
	fi

	# Defensive: if a previous start_service installed nft rules and the
	# user has since flipped every tproxy inbound to a non-tproxy protocol
	# (or disabled the section), `needed` will now return "0" — but the
	# old table would linger until the next `stop`. Remove unconditionally
	# first so the post-apply state always reflects the latest `needed`
	# decision. The remove is idempotent and a no-op when no table exists.
	ucode -L "$UCODE_LIB" "$NFTABLES_UC" remove 2>/dev/null || true

	if [ "$(ucode -L "$UCODE_LIB" "$NFTABLES_UC" needed)" = "1" ]; then
		# G4: surface apply failures via syslog. The script's atomic
		# `add+delete+table` prelude means the prior live table is
		# untouched on failure, so we don't try to roll back here —
		# we only ensure the operator sees the rc.
		ucode -L "$UCODE_LIB" "$NFTABLES_UC" apply
		_nft_rc=$?
		[ "$_nft_rc" = 0 ] || logger -t singbox-ui "nft apply failed (rc=$_nft_rc)"
	fi

	procd_open_instance
	# C3.6: drop to minimal capability set so sing-box runs without root.
	# We intentionally do NOT switch user/group to nobody/nogroup here —
	# that requires fixing FS-permissions on /tmp/singbox-ui.json and the
	# /tmp/singbox-ui/ subscription cache, which could break test envs.
	# Capabilities-only is a real security win without that risk.
	procd_set_param capabilities /etc/capabilities/singbox-ui.json
	procd_set_param command "$SINGBOX_BIN" run -c /tmp/singbox-ui.json
	procd_set_param respawn
	procd_close_instance

	# C3.3: record service-start timestamp for status_detail RPC. S1.1: the
	# config passed `sing-box check` above and was handed to procd, so record the
	# real apply result here — status_detail surfaces this as last_apply_result,
	# distinct from generate.uc's last_generate_result.
	mkdir -p /var/lib/singbox-ui
	echo "{\"service_start_ts\":$(date +%s)}" > /var/lib/singbox-ui/service_state.json
	echo "{\"last_apply_ts\":$(date +%s),\"last_apply_result\":\"ok\"}" \
		> /var/lib/singbox-ui/apply_state.json

	# Cache the package version once per start so status_detail (a hot read
	# RPC) doesn't fork apk on every call. Refreshed here so an upgrade +
	# restart re-derives it. (spec S5-PERF). OpenWrt 25.12.3 ships apk, not
	# opkg; `apk list --installed` prints "<name>-<version> noarch {...}".
	# Must yield the SAME cache format as cached_pkg_version in the rpcd
	# handler: the version with the package-name prefix stripped.
	apk list --installed luci-singbox-ui 2>/dev/null \
		| sed -n 's/^luci-singbox-ui-\([^ ]*\) .*/\1/p' \
		> /var/lib/singbox-ui/pkg_version || true
}

stop_service() {
	# G7: nftables.uc remove prints "No such file or directory" on a
	# fresh stop (no table to delete). Symmetric with start_service
	# (line above already suppresses); keep operator logs quiet.
	ucode -L "$UCODE_LIB" "$NFTABLES_UC" remove 2>/dev/null
}

reload_config() {
	# sing-box does not support signal-based config reload.
	stop
	start
}
